# deploy/systemd/diffdelta-agent.service
# DiffDelta Agent — systemd service for Linux servers.
#
# Setup (3 commands):
#   1. pip install diffdelta
#   2. Copy a recipe to /opt/diffdelta/:
#        sudo mkdir -p /opt/diffdelta
#        sudo cp recipes/slack_alerts.py /opt/diffdelta/agent.py
#        sudo chown -R nobody:nogroup /opt/diffdelta
#   3. Install this service:
#        sudo cp diffdelta-agent.service /etc/systemd/system/
#        sudo systemctl daemon-reload
#        sudo systemctl enable --now diffdelta-agent
#
# Logs:
#   journalctl -u diffdelta-agent -f

[Unit]
Description=DiffDelta Intelligence Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=nobody
Group=nogroup

# Cursor persistence — saved to /opt/diffdelta/cursors.json
# so it survives restarts. (Default ~/ path won't work for nobody user.)
Environment="DD_CURSOR_PATH=/opt/diffdelta/cursors.json"

# Set your integration webhook/token here:
Environment="SLACK_WEBHOOK_URL=https://hooks.slack.com/services/T.../B.../xxx"
# Environment="DISCORD_WEBHOOK_URL=https://discord.com/api/webhooks/..."
# Environment="GITHUB_TOKEN=ghp_..."
# Environment="GITHUB_REPO=your-org/your-repo"
# Environment="PAGERDUTY_ROUTING_KEY=R..."

ExecStart=/usr/bin/python3 /opt/diffdelta/agent.py

Restart=always
RestartSec=30

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/diffdelta

[Install]
WantedBy=multi-user.target
